# indexer/crontab
@daily python $HOME/index.py &>> $STACKATO_APP_ROOT/logs/cron.log
